# luci-app-xwan/files/root/etc/uci-defaults/40_luci-app-xwan

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
