# luci-app-turboacc/root/etc/uci-defaults/luci-turboacc
#!/bin/sh

uci -q batch <<-EOF >/dev/null
EOF

[ -n "$(lsmod | grep mtkhnat)" ] && \
	uci set turboacc.config.sw_flow='0' && \
	uci set turboacc.config.hw_flow='0'

uci commit turboacc

exit 0
