# luci-app-cloudreve/root/etc/uci-defaults/cloudreve
#!/bin/sh

uci -q batch <<-EOF >/dev/null
EOF

chmod +x /etc/init.d/cloudreve >/dev/null 2>&1

rm -f /tmp/luci-indexcache
exit 0
