# luci-app-autotimeset/root/etc/uci-defaults/luci-autotimeset
#!/bin/sh
chmod +x  /etc/init.d/autotimeset /usr/bin/timesethandler /etc/init.d/autotimesetrun
uci -q batch <<-EOF >/dev/null
EOF
rm -rf /tmp/luci-modulecache /tmp/luci-indexcache*
exit 0