# luci-app-dnspod/root/etc/uci-defaults/luci-dnspod
#!/bin/sh

uci -q batch <<-EOF >/dev/null
EOF


/etc/init.d/dnspod stop

rm -f /tmp/luci-indexcache
exit 0
