# luci-app-ipsec-ikev2/root/etc/init.d/ikev2
#!/bin/sh /etc/rc.common
# IKEv2 LuCI helper - generates swanctl.conf and reloads strongSwan (PSK)
START=99
USE_PROCD=0

_ENABLED="$(uci -q get ikev2.config.enabled)"
GEN="/usr/libexec/ikev2/genconfig.sh"

start() {
    [ "$_ENABLED" = "1" ] || {
        echo "luci-ikev2: disabled"
        return 0
    }
    [ -x "$GEN" ] || {
        echo "luci-ikev2: generator not found: $GEN"
        return 1
    }
    "$GEN" || return 1
    /etc/init.d/strongswan enable >/dev/null 2>&1
    /etc/init.d/strongswan start
    swanctl --load-all >/dev/null 2>&1
    echo "luci-ikev2: started"
}

stop() {
    /etc/init.d/strongswan stop
    echo "luci-ikev2: stopped"
}

reload() {
    [ "$_ENABLED" = "1" ] || {
        /etc/init.d/strongswan stop
        return 0
    }
    "$GEN" || return 1
    /etc/init.d/strongswan start
    swanctl --load-all >/dev/null 2>&1
    /etc/init.d/firewall reload >/dev/null 2>&1
    echo "luci-ikev2: reloaded"
}
