# luci-app-lucky/root/etc/uci-defaults/luci-lucky
#!/bin/sh

uci -q batch <<-EOF >/dev/null
EOF
[[ -f /etc/config/lucky ]] || echo "config lucky 'lucky' " > /etc/config/lucky
rm -rf /tmp/luci*
exit 0
