# luci-app-qos-emong/files/etc/uci-defaults/luci-qos-emong
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
EOF
rm -f /tmp/luci-indexcahe
exit 0
