# luci-app-xlnetacc/root/etc/uci-defaults/luci-xlnetacc
uci batch <<EOF
EOF
rm -rf /tmp/luci-*
exit 0
