# .github/diy/packages/luci-app-samba4/root/etc/uci-defaults/luci-samba4
#!/bin/sh

chmod +x /etc/hotplug.d/block/20-smb

exit 0
