# luci-app-fan/root/etc/uci-defaults/luci-fan
#!/bin/sh

/usr/libexec/fan-control init

uci -q batch <<-EOF >/dev/null
EOF

exit 0
